# backend/docker/supervisord.conf
; ═══════════════════════════════════════════════════════════════
;  Supervisor — runs both nginx + php-fpm inside the backend container.
;  If either process dies, supervisor restarts it. If supervisor itself
;  exits non-zero, Docker restarts the whole container.
; ═══════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true                          ; run in foreground (PID 1 inside container)
user=root                              ; supervisord drops privileges via the
                                       ; programs' own `user=` directives below
logfile=/dev/null                      ; logs go to programs' own stdout/stderr
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10                            ; start before nginx
stdout_logfile=/dev/fd/1               ; → docker logs
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
user=root                              ; FPM master needs root to bind socket;
                                       ; worker processes run as `panelos` per www.conf

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20                            ; start after php-fpm has the socket
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
user=root                              ; nginx master needs root to bind :8080;
                                       ; worker processes auto-drop to `nginx`
